# Chihiro/docker/chilhrio-production-environment/monitor.ini
[program:monitor]
command=/usr/bin/python3 monitor.py
user = root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
startretries=5
directory=/Chihiro
stopsignal=KILL
priority=50