# Chihiro/docker/chilhrio-production-environment/v2ray.ini
[program:v2ray]
command=/usr/bin/v2ray/v2ray -config=/Chihiro/config.json
user = root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
startretries=9999
directory=/usr/bin/v2ray
stopsignal=KILL
priority=100